# example/autostart/ubuntu/etc/systemd/system/shhoook.service
[Unit]
Description=shhoook — impress the server (HTTP → shell hooks)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/shhoook
ExecStart=/usr/local/bin/shhoook-wrapper
Restart=on-failure
RestartSec=2

# Minimal hardening (shhoook executes scripts, so keep it conservative)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
