# example/autostart/ubuntu-iface/etc/systemd/system/shhoook.service
[Unit]
Description=shhoook — impress the server (HTTP → shell hooks)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# You can override these in /etc/default/shhoook
EnvironmentFile=-/etc/default/shhoook

# Defaults (if /etc/default/shhoook is missing)
Environment=CONFIG_DIR=/etc/shhoook/conf
Environment=WAIT_IFACE=tun0
Environment=PORT=8080
Environment=WAIT_IP=
Environment=WAIT_V4_PREFIX=

# Compute LISTEN_ADDR from interface IPv4 (optionally wait for a specific IP or prefix),
# then exec the daemon with LISTEN_ADDR exported.
ExecStart=/bin/sh -lc '\
  set -eu; \
  iface="${WAIT_IFACE}"; port="${PORT}"; want_ip="${WAIT_IP}"; prefix="${WAIT_V4_PREFIX}"; \
  i=0; \
  while [ "$i" -lt 60 ]; do \
    if ip link show dev "$iface" >/dev/null 2>&1; then \
      # Get all IPv4 addresses on the interface (space-separated)
      ips="$(ip -4 -o addr show dev "$iface" | awk "{print \$4}" | cut -d/ -f1 | xargs echo)"; \
      picked=""; \
      if [ -n "$want_ip" ]; then \
        echo " $ips " | grep -q " $want_ip " && picked="$want_ip" || picked=""; \
      elif [ -n "$prefix" ]; then \
        for x in $ips; do case "$x" in ${prefix}*) picked="$x"; break;; esac; done; \
      else \
        for x in $ips; do picked="$x"; break; done; \
      fi; \
      if [ -n "$picked" ]; then \
        export LISTEN_ADDR="${picked}:${port}"; \
        export CONFIG_DIR="${CONFIG_DIR}"; \
        echo "Starting shhoook with LISTEN_ADDR=${LISTEN_ADDR} (from ${iface})"; \
        exec /usr/local/bin/shhoook; \
      fi; \
    fi; \
    i=$((i+1)); sleep 1; \
  done; \
  echo "Interface not ready: ${iface} (wanted IP: ${want_ip:-any}, prefix: ${prefix:-any})" >&2; \
  exit 1'

Restart=on-failure
RestartSec=2

# Reasonable hardening (keep it minimal; shhoook executes scripts)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target